# deploy/options-bidbot.service
[Unit]
Description=BidBot Options Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bidbot
Environment=PATH=/opt/bidbot/.venv/bin
ExecStart=/opt/bidbot/.venv/bin/python scripts/run_strategy.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
